# deploy/supervisor/elasticsearch-prod.conf
; Supervisor process profile for Nexus Elasticsearch production instance.
; Copy or symlink deploy/conf/elasticsearch-prod.conf to
; /etc/elasticsearch/elasticsearch.yml before starting this program.

[program:nexus-elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
directory=/usr/share/elasticsearch
priority=180
autostart=true
autorestart=unexpected
startsecs=45
startretries=5
exitcodes=0
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
umask=0077
stdout_logfile=/var/log/nexus/elasticsearch-supervisor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nexus/elasticsearch-supervisor.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
environment=ES_PATH_CONF="/etc/elasticsearch",ES_JAVA_OPTS="-Xms4g -Xmx4g",ES_NODE_NAME="es-1"

[program:nexus-elasticsearch-watchdog]
command=/opt/nexus/deploy/supervisor/middleware-watchdog.sh nexus-elasticsearch /opt/nexus/deploy/supervisor/healthcheck-elasticsearch.sh
user=elasticsearch
priority=190
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/nexus/elasticsearch-watchdog.log
stderr_logfile=/var/log/nexus/elasticsearch-watchdog.err.log
environment=CHECK_INTERVAL="30",MAX_FAILURES="3",RESTART_COOLDOWN="180",ES_HEALTH_SCHEME="https",ES_HEALTH_HOST="127.0.0.1",ES_HEALTH_PORT="9200",ES_USERNAME="elastic",ES_PASSWORD_FILE="/run/secrets/elasticsearch_health_password"
