# deploy/supervisor/minio-prod.conf
; Supervisor process profile for Nexus MinIO production instance.
; MinIO high availability comes from distributed erasure coding across nodes.
; Supervisor only keeps the local server process alive.

[program:nexus-minio]
command=/bin/bash -lc 'set -a; source /opt/nexus/deploy/conf/minio-prod.conf; set +a; exec /usr/local/bin/minio server ${MINIO_OPTS} ${MINIO_VOLUMES}'
user=minio
directory=/var/lib/minio
priority=140
autostart=true
autorestart=unexpected
startsecs=20
startretries=5
exitcodes=0
stopsignal=TERM
stopwaitsecs=60
stopasgroup=true
killasgroup=true
umask=0077
stdout_logfile=/var/log/nexus/minio-supervisor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nexus/minio-supervisor.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10

[program:nexus-minio-watchdog]
command=/opt/nexus/deploy/supervisor/middleware-watchdog.sh nexus-minio /opt/nexus/deploy/supervisor/healthcheck-minio.sh
user=minio
priority=150
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/nexus/minio-watchdog.log
stderr_logfile=/var/log/nexus/minio-watchdog.err.log
environment=CHECK_INTERVAL="30",MAX_FAILURES="3",RESTART_COOLDOWN="120",MINIO_HEALTH_SCHEME="https",MINIO_HEALTH_HOST="127.0.0.1",MINIO_HEALTH_PORT="9000"
