# deploy/supervisor/postgresql-prod.conf
; Supervisor process profile for Nexus PostgreSQL production instance.
; Process supervision is not a substitute for database HA. Use streaming
; replication/Patroni/cloud-managed HA when RTO/RPO requires automatic failover.

[program:nexus-postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/opt/nexus/deploy/conf/postgresql-prod.conf
user=postgres
directory=/var/lib/postgresql
priority=100
autostart=true
autorestart=unexpected
startsecs=15
startretries=5
exitcodes=0
stopsignal=INT
stopwaitsecs=90
stopasgroup=true
killasgroup=true
umask=0077
stdout_logfile=/var/log/nexus/postgresql-supervisor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nexus/postgresql-supervisor.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
environment=PGDATA="/var/lib/postgresql/15/main"

[program:nexus-postgresql-watchdog]
command=/opt/nexus/deploy/supervisor/middleware-watchdog.sh nexus-postgresql /opt/nexus/deploy/supervisor/healthcheck-postgresql.sh
user=postgres
priority=110
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/nexus/postgresql-watchdog.log
stderr_logfile=/var/log/nexus/postgresql-watchdog.err.log
environment=CHECK_INTERVAL="30",MAX_FAILURES="3",RESTART_COOLDOWN="120",PGHOST="127.0.0.1",PGPORT="5432",PGUSER="postgres"
