# deploy/supervisor/redis-prod.conf
; Supervisor process profile for Nexus Redis production instance.
; For cache-critical production, prefer Redis Sentinel/Cluster in addition to
; this process-level watchdog.

[program:nexus-redis]
command=/usr/bin/redis-server /opt/nexus/deploy/conf/redis-prod.conf --daemonize no --supervised no
user=redis
directory=/var/lib/redis
priority=120
autostart=true
autorestart=unexpected
startsecs=10
startretries=5
exitcodes=0
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
umask=0077
stdout_logfile=/var/log/nexus/redis-supervisor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nexus/redis-supervisor.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10

[program:nexus-redis-watchdog]
command=/opt/nexus/deploy/supervisor/middleware-watchdog.sh nexus-redis /opt/nexus/deploy/supervisor/healthcheck-redis.sh
user=redis
priority=130
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/var/log/nexus/redis-watchdog.log
stderr_logfile=/var/log/nexus/redis-watchdog.err.log
environment=CHECK_INTERVAL="30",MAX_FAILURES="3",RESTART_COOLDOWN="120",REDIS_HOST="127.0.0.1",REDIS_PORT="6379",REDIS_PASSWORD_FILE="/run/secrets/redis_health_password"
